# infra/ansible/roles/cd_deploy/templates/devsecops-deploy.timer.j2
[Unit]
Description=Verifie main et deploie periodiquement le DevSecOps lab

[Timer]
OnBootSec=2min
OnUnitActiveSec={{ cd_deploy_interval }}
Persistent=true

[Install]
WantedBy=timers.target
